# miscellaneous/lb.service
[Unit]
Description=Look busy (at 15%)

[Service]
ExecStart=/usr/local/bin/lb -c 15

[Install]
WantedBy=multi-user.target